# mongo-db/7.0.2/conf/supervisor_compose.ini
[program:mongo-db]
command=bash -c "${MONGO_DB_RUN_COMMAND}"
process_name=%(program_name)s
pidfile=/var/run/mongod.pid
autorestart=false
