# workload/xen-orchestra/6/conf/supervisor.ini
[program:xen-orchestra]
command=bash -c "yarn start"
process_name=%(program_name)s
directory=/usr/local/lib/xen-orchestra/packages/xo-server
pidfile=/opt/xen-orchestra/system/server.pid
autorestart=false
