# workload/kafka/3.2.1/controller/conf/supervisor.ini
[program:kafka]
command=bash -c "/opt/kafka/script/start.sh"
process_name=%(program_name)s
pidfile=/opt/kafka/system/node.pid
numprocs=1
autorestart=false
