# workload/mongo-db/4.4.4/replica-cluster/conf/supervisor.ini
[group:mongo-db]
programs=init,server

[program:init]
command=bash -c "mongosh /opt/mongo-db/script/init.js --port 7000"
process_name=%(program_name)s
autorestart=false
exitcodes=0,1
priority=100

[program:server]
command=bash -c "sleep 3 && mongod -f /opt/mongo-db/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/mongo-db/system/server.pid
autorestart=false
priority=200
