# workload/grafana/12.3.0/conf/supervisor.ini
[program:grafana]
command=bash -c "/usr/local/lib/grafana/bin/grafana server --config /opt/grafana/conf/server.ini"
process_name=%(program_name)s
pidfile=/opt/grafana/system/server.pid
autorestart=false
