# workload/grafana/12.3.0/conf/supervisor.ini
[program:grafana]
command=bash -c "/usr/local/bin/grafana/bin/grafana server --homepath /usr/local/bin/grafana"
process_name=%(program_name)s
pidfile=/opt/grafana/system/server.pid
autorestart=false
