# workload/postgresql/16.4/conf/supervisor.ini
[program:postgresql]
command=bash -c "postgres -c /opt/postgresql/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/mongo-db/system/server.pid
autorestart=false
