# workload/postgresql/16.4/conf/supervisor.ini
[group:postgresql]
programs=init,server

[program:init]
command=su postgresql bash -c "initdb -D /opt/postgresql/data"
process_name=%(program_name)s
autorestart=false
exitcodes=0,1
priority=100

[program:server]
command=su postgresql bash -c "sleep 3 && postgres -c config_file=/opt/postgresql/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/postgresql/system/server.pid
autorestart=false
priority=200
