# workload/chroma/1.1.4/conf/supervisor.ini
[program:chroma]
command=bash -c "chroma run /opt/chroma/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/chroma/system/server.pid
autorestart=false
