# workload/kibana/8.17.0/conf/supervisor.ini
[program:kibana]
command=bash -c "/usr/local/lib/kibana/bin/kibana -c /opt/kibana/conf/server.yml --allow-root"
process_name=%(program_name)s
pidfile=/opt/kibana/system/server.pid
autorestart=false
