# workload/prometheus/2.53.1/conf/supervisor.ini
[program:prometheus]
command=bash -c "prometheus --config.file=/opt/prometheus/conf/prometheus.yml --storage.tsdb.path=/opt/prometheus/data --web.console.templates=/opt/prometheus/consoles --web.console.libraries=/opt/prometheus/console_libraries"
process_name=%(program_name)s
pidfile=/opt/prometheus/system/server.pid
autorestart=false