# mongo-db/4.4.4/conf/supervisor.ini
[program:mongo-db]
command=bash -c "mongod -f /opt/mongo-db/conf/mongod.conf"
process_name=%(program_name)s
pidfile=/var/run/mongod.pid
autorestart=false
