# _stack/ctm/app/conf/supervisor.ini
[group:app]
programs=frontend,backend

[program:frontend]
command=bash -c "npm start --port=3000"
process_name=%(program_name)s
directory=/opt/app/source/frontend
stdout_logfile=/opt/app/log/frontend.log
stderr_logfile=/opt/app/log/frontend-error.log
autorestart=false
exitcodes=0,1
priority=200

[program:backend]
command=bash -c "npm run start --port=2000"
process_name=%(program_name)s
directory=/opt/app/source/backend
stdout_logfile=/opt/app/log/backend.log
stderr_logfile=/opt/app/log/backend-error.log
autorestart=true
priority=100
