# workload/elasticsearch/9.2.3/conf/supervisor.ini
[program:elasticsearch]
command=/usr/local/lib/elasticsearch/bin/elasticsearch -p /opt/elasticsearch/system/server.pid
user=elasticsearch
directory=/usr/local/lib/elasticsearch
process_name=%(program_name)s
pidfile=/opt/elasticsearch/system/server.pid
autorestart=false
