# workload/apm/6.8.2/conf/supervisor.ini
[program:apm]
command=bash -c "/usr/local/lib/apm/apm-server run -e --path.config=/usr/local/lib/apm --strict.perms=false"
process_name=%(program_name)s
pidfile=/opt/apm/system/server.pid
autorestart=false
