# supervisor/conf/app.ini
[program:app]
command=/bin/bash -c "python init.py"
process_name=%(program_name)s
directory=/opt/app
numprocs=1
autorestart=false