# supervisor/conf/app.ini
[program:app]
command=/bin/bash -c "python init.py"
process_name=%(program_name)s
directory=/opt/app
autorestart=false
stdout_logfile=/opt/app/log/stdout.log
stderr_logfile=/opt/app/log/stderr.log
