# supervisor/conf/redis.ini
[program:redis]
command=bash -c "/usr/local/bin/redis-server /opt/redis/conf/redis.conf"
process_name=%(program_name)s
pidfile=/var/run/redis.pid
autorestart=false
