# workload/logstash/9.2.3/conf/supervisor.ini
[program:logstash]
command=/usr/local/lib/logstash/bin/logstash --path.settings /opt/logstash/conf
user=logstash
directory=/usr/local/lib/logstash
process_name=%(program_name)s
pidfile=/opt/logstash/system/server.pid
autorestart=false