# workload/elasticsearch/8.17.0/conf/supervisor.ini
[program:elasticsearch]
command=su elasticsearch bash -c "/usr/local/lib/elasticsearch/bin/elasticsearch -p /opt/elasticsearch/system/server.pid"
process_name=%(program_name)s
pidfile=/opt/elasticsearch/system/server.pid
autorestart=false
