# workload/hadoop/2.10.1/name/conf/supervisor.ini
[program:hadoop]
command=bash -c "/opt/hadoop/script/start.sh"
process_name=%(program_name)s
pidfile=/opt/hadoop/system/node.pid
numprocs=1
autorestart=false
