# nginx/1.14.0/conf/supervisor.ini
[program:nginx]
command=bash -c "nginx"
process_name=%(program_name)s
pidfile=/var/run/nginx.pid
autorestart=false
