# redis/4.0.9/conf/supervisor.ini
[program:redis]
command=bash -c "redis-server /opt/redis/conf/redis.conf"
process_name=%(program_name)s
pidfile=/var/run/redis.pid
autorestart=false
