# workload/mongo-db/4.4.4/conf/supervisor.ini
[program:mongo-db]
command=bash -c "mongod --replSet replica_cluster --config /opt/mongo-db/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/mongo-db/system/server.pid
autorestart=true

[program:init-replica]
command=bash -c "chmod +x /opt/mongo-db/script/init-replica-set.sh && /opt/mongo-db/script/init-replica-set.sh"
process_name=%(program_name)s
autorestart=false     ; Disable auto-restart since the script runs only once
