# workload/mongo-db/4.4.4/conf/supervisor.ini
[program:mongo-db]
command=bash -c "mongod -f /opt/mongo-db/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/mongo-db/system/server.pid
autorestart=false
