# _stack/ntm/app/conf/supervisor.ini
[program:app]
command=bash -c "rails server -b 0.0.0.0 -p 2000"
process_name=%(program_name)s
directory=/opt/app/source
autorestart=false
stdout_logfile=/opt/app/log/stdout.log
stderr_logfile=/opt/app/log/stderr.log
