# workload/nginx/1.14.0/conf/supervisor.ini
[program:nginx]
command=bash -c "nginx -c /opt/nginx/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/nginx/system/server.pid
autorestart=false
