# workload/redis/7.2.1/conf/supervisor.ini
[program:redis]
command=bash -c "redis-server /opt/redis/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/redis/system/server.pid
autorestart=false
