# workload/prometheus/3.8.1/conf/supervisor.ini
[program:prometheus]
command=bash -c "prometheus --config.file=/opt/prometheus/conf/server.conf"
process_name=%(program_name)s
pidfile=/opt/prometheus/system/server.pid
autorestart=false
