# workload/ollama/0.9.6/conf/supervisor.ini
[program:ollama]
command=bash -c "ollama serve"
process_name=%(program_name)s
pidfile=/opt/ollama/system/server.pid
autorestart=false
