# workload/ollama/0.9.6/gemma/2/conf/supervisor.ini
[group:ollama]
programs=server,pull

[program:server]
command=bash -c "ollama serve"
process_name=%(program_name)s
pidfile=/opt/ollama/system/server.pid
autorestart=false
exitcodes=0,1
priority=100

[program:pull]
command=bash -c "sleep 3 && ollama pull gemma2:2b"
process_name=%(program_name)s
pidfile=/opt/ollama/system/server.pid
autorestart=false
autostart=false
priority=200
