# workload/apm/9.2.3/conf/supervisor.ini
[program:apm]
command=/usr/local/lib/apm/apm-server run --path.config=/usr/local/lib/apm --strict.perms=false
directory=/usr/local/lib/apm
process_name=%(program_name)s
pidfile=/opt/apm/system/server.pid
autorestart=false
