# workload/logstash/8.17.0/conf/supervisor.ini
[program:logstash]
command=bash -c "/usr/local/lib/logstash/bin/logstash"
process_name=%(program_name)s
pidfile=/opt/logstash/system/server.pid
autorestart=false
